# systemd/.config/systemd/user/app-mako@autostart.service.d/restart.conf
# Keep the notification daemon alive.
#
# The XDG-autostart-generated unit ships Restart=no, so once mako dies it stays
# dead until the next theme switch or a manual start. On 2026-06-12 a theme-switch
# restart storm (10+ `systemctl --user restart app-mako@autostart` in under a
# minute) tripped systemd's default start-rate-limit and wedged the unit; mako was
# down ~21:00–21:44, and the evening report crons' notify-send failed (the packaged
# fallback mako.service is masked). See household-oc memory project_cron_notify_reliability.
#
# StartLimitIntervalSec=0 → the restart storm can't wedge the unit.
# Restart=on-failure       → a crashed daemon comes back on its own.
[Unit]
StartLimitIntervalSec=0

[Service]
Restart=on-failure
RestartSec=2
