# autostart/.config/autostart/hypridle.desktop
[Desktop Entry]
Type=Application
Name=Hypridle
Exec=hypridle
X-GNOME-Autostart-enabled=true
