# autostart/.config/autostart/hyprpaper.desktop
[Desktop Entry]
Type=Application
Name=Hyprpaper
Exec=hyprpaper
X-GNOME-Autostart-enabled=true
