# autostart/.config/autostart/mako.desktop
[Desktop Entry]
Type=Application
Name=Mako
Exec=mako
X-GNOME-Autostart-enabled=true
