# autostart/.config/autostart/theme-init.desktop
[Desktop Entry]
Type=Application
Name=Theme Initializer
Comment=Apply saved theme on startup
Exec=/home/curator/.dotfiles/bin/.bin/theme-startup.sh
Terminal=false
X-GNOME-Autostart-enabled=true
