# autostart/.config/autostart/theme-init.desktop
[Desktop Entry]
Type=Application
Name=Theme Initializer
Comment=Pick a random theme once per boot, then keep it stable
Exec=/home/curator/.dotfiles/bin/.bin/theme-startup.sh
Terminal=false
X-GNOME-Autostart-enabled=true
