# autostart/.config/autostart/waybar.desktop
[Desktop Entry]
Type=Application
Name=Waybar
Exec=waybar
X-GNOME-Autostart-enabled=true
