# autostart/.config/autostart/wl-clip-persist.desktop
[Desktop Entry]
Type=Application
Name=wl-clip-persist
Exec=wl-clip-persist --clipboard both
X-GNOME-Autostart-enabled=true
